# .dist/linux/etc/profile.d/bearded-wookie-completion.sh
#!/bin/sh

eval "$(bw install-completions)"