# .test/linux-development/systemd/bearded-wookie-pebble.service
[Unit]
Description=bearded-wookie agent
After=network.target

[Service]
Type=simple
Restart=always
StateDirectory=%p/%i
RuntimeDirectory=%p/%i
CacheDirectory=%p/%i
ConfigurationDirectory=%p/%i
WorkingDirectory=%S/%p
SyslogIdentifier=%p
Environment=PEBBLE_VA_NOSLEEP=1
Environment=PEBBLE_WFE_NONCEREJECT=0
ExecStartPre=/usr/bin/echo "cache=${CACHE_DIRECTORY} state=${STATE_DIRECTORY} runtime=${RUNTIME_DIRECTORY} - %S/%p/config.json"
ExecStartPre=/usr/bin/tree %S
ExecStartPre=/usr/bin/echo "%h/go/bin/pebble -strict -config %S/%p/config.json"
ExecStart=%h/go/bin/pebble -strict -config %S/%p/config.json

[Install]
WantedBy=default.target
