# .dist/linux/usr/lib/systemd/system/bearded-wookie.slice
[Unit]
Description=bearded-wookie resource limits; primarily to make it easy to protect systems it runs on.
DefaultDependencies=no
Before=slices.target

[Slice]
CPUAccounting=true
CPUQuota=25%